# systemd/udot-api.service
[Unit]
Description=UDOT Road Conditions API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/utah-road-conditions-api

# Use virtual environment's gunicorn and Python
ExecStart=/root/utah-road-conditions-api/venv/bin/gunicorn -w 4 -b 0.0.0.0:5000 api_server:app

# Restart on failure
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/udot-api.log
StandardError=append:/var/log/udot-api-error.log

[Install]
WantedBy=multi-user.target
